# deploy/sync-analyzer-celery.service
[Unit]
Description=Professional Audio Sync Analyzer - Celery Worker
Documentation=https://github.com/your-org/sync-analyzer
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=sync-analyzer
Group=sync-analyzer
WorkingDirectory=/opt/sync-analyzer/fastapi_app

# Environment
Environment="PYTHONPATH=/opt/sync-analyzer:/opt/sync-analyzer/fastapi_app"
Environment="PYTHONUNBUFFERED=1"
Environment="REDIS_URL=redis://localhost:6379/0"
EnvironmentFile=-/opt/sync-analyzer/.env

# Exec
ExecStart=/opt/sync-analyzer/venv/bin/celery -A app.core.celery_app worker --loglevel=info --concurrency=4 --hostname=sync-worker@%%h
ExecReload=/bin/kill -HUP $MAINPID

# Graceful shutdown
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sync-analyzer/logs /opt/sync-analyzer/uploads /opt/sync-analyzer/reports /mnt/data
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=append:/opt/sync-analyzer/logs/celery.log
StandardError=append:/opt/sync-analyzer/logs/celery_error.log

[Install]
WantedBy=multi-user.target

